# linux/pingwatch.service
[Unit]
Description=PingWatch Network Monitor
After=network.target postgresql.service
Wants=postgresql.service
# Bound crash loops: if startup fails 5 times within 120s (e.g. a bad deploy
# with a SyntaxError), give up and enter 'failed' instead of restarting every
# 5s forever. Previously 0 (unlimited), which let a single bad pull loop
# indefinitely and flood the logs with startup banners.
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/pingwatch
# Refuse to start if any source file has a syntax error. compileall is purely
# syntactic (no imports, no side effects), so a broken deploy fails fast here
# and counts toward StartLimitBurst instead of half-starting. venv is excluded.
ExecStartPre=/usr/bin/python3 -m compileall -q -x venv /opt/pingwatch
ExecStart=/usr/bin/python3 /opt/pingwatch/server.py
Restart=on-failure
RestartSec=5

# Run as the installing user — patched automatically by start.sh --install-service
# User=pingwatch
# Group=pingwatch

# Allow binding to privileged ports (80, 443, 162) without full root
# CAP_NET_RAW   — required for ICMP ping (raw socket)
# CAP_NET_BIND_SERVICE — required for ports < 1024
AmbientCapabilities=CAP_NET_RAW CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_BIND_SERVICE

# Suppress console output — logs go to journald (journalctl -u pingwatch)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pingwatch

[Install]
WantedBy=multi-user.target
